# systemd/tracer.timer
[Unit]
Description=Runs every 30 min tracer.py

[Timer]
OnCalendar=*:0/30

[Install]
WantedBy=timers.target