# systemd/tracer.timer
[Unit]
Description=Runs every 60 min tracer.py

[Timer]
OnCalendar=00/1:00

[Install]
WantedBy=timers.target